# health_report.service
[Unit]
Description=Daily Health Report

[Service]
Type=oneshot
User=soufiane
ExecStart=/bin/bash /home/soufiane/linux-health-monitor/report.sh
